# service/sysmon.service
[Unit]
Description=System monitoring and configuration checking service

[Service]
Type=oneshot
ExecStart=/usr/bin/sysmon
