# service/system-monitor.service
[Unit]
Description=System monitoring and configuration checking service

[Service]
Type=oneshot
ExecStart=/opt/system-monitor/bootstrap.sh
